# deploy/systemd/tiku-zong-api.service
[Unit]
Description=TIKU-ZONG Go API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tiku
Group=tiku
WorkingDirectory=/opt/tiku-zong/backend
EnvironmentFile=/etc/tiku-zong/tiku-zong.env
ExecStart=/opt/tiku-zong/backend/tiku-zong-api
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
